# build/microk8s-systemd/mdns-ingress-publisher.service
[Unit]
Description=mdns-ingress-publisher - Watches for Kubernetes Ingress object events and publishes mDNS records accordingly
After=network.target

[Service]
User=mdns-ingress-publisher
Type=simple
WorkingDirectory=/home/mdns-ingress-publisher
ExecStart=mdns-ingress-publisher
Restart=on-failure
PrivateTmp=true
PrivateDevices=true

[Install]
WantedBy=multi-user.target